# template/sysstart_aes67.service.j2
[Unit]
Description=Correct running daemon AES67
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/local/run_aes67.sh

[Install]
WantedBy=multi-user.target